# pgxl/etc/supervisor/conf.d/pgxl.conf
[program:pgxl]
user=pom
command=/app/bin/pgxl.sh
redirect_stderr=true
stopsignal=QUIT
