# pgxl/etc/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[inet_http_server]
port=9001

[program:sshd]
priority=100
command=sudo /usr/sbin/sshd -D
autostart=true
autorestart=true
redirect_stderr=true
